# automation/antigravity-remote.service
[Unit]
Description=Antigravity Mobile Remote Service
After=network.target

[Service]
User=hawkatoad
WorkingDirectory=/home/hawkatoad/Desktop/anti gravity/passive-income-store/automation
ExecStart=/usr/bin/python3 "/home/hawkatoad/Desktop/anti gravity/passive-income-store/automation/mobile_remote.py"
Restart=always

[Install]
WantedBy=multi-user.target
